# autoscaler.service
[Unit]
Description=Docker Swarm Auto Scaler
# Docker가 실행된 후에 시작
After=docker.service
Requires=docker.service

[Service]
Type=simple
# Manager 노드의 ubuntu 계정으로 실행
User=ubuntu
ExecStart=/usr/local/bin/autoscale.sh
# 비정상 종료 시 자동 재시작
Restart=always
RestartSec=5s
# 로그를 journald에도 기록
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
